# gentoo/init.d/shadowsocks-rust
#!/sbin/openrc-run

depend() {
    after net
}

description="a fast tunnel proxy that helps you bypass firewalls"
command="/usr/local/bin/ssserver"
command_user="zkc:zkc"
command_args="-c /usr/local/etc/shadowsocks-rust/config.json"
supervisor="supervise-daemon"
