# gentoo/init.d/skb
#!/sbin/openrc-run

depend() {
    after net
}

description="simple knownledge base"
command="/home/zkc/projects/skb/bin/skb"
command_user="zkc:zkc"
supervisor="supervise-daemon"
supervise_daemon_args="-d /home/zkc/projects/skb"
